# systemd/.config/systemd/user/syncthing.service.d/override.conf
[Service]
ExecStartPre=/bin/bash -c 'while [ "$(curl http://www.archlinux.org/check_network_status.txt)" != "NetworkManager is online" ]; do sleep 1; done'
